# supervisord.conf
[supervisord]
nodaemon=true

[program:topology]
command=ipython /code/flask/app.py