# supervisord.conf
[supervisord]
nodaemon=true

[program:topology]
command=python /code/flask/app.py